# linux/ecosystem/torrust-index/supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
killasgroup=true
stopasgroup=true
redirect_stderr=true

[program:torrust-index-backend]
directory=/app/bin
command=/bin/bash -c "/app/bin/torrust-index-backend"
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log