# linux/ecosystem/teamcity/agent/minimal/etc/supervisor/conf.d/teamcity-agent.conf
[program:teamcity-agent]
command=/usr/local/share/epicmorg/teamcity/agent/sbin/entrypoint.sh
environment=JAVA_HOME="/usr/local/share/epicmorg/java/21",JRE_HOME="/usr/local/share/epicmorg/java/21",JDK_HOME="/usr/local/share/epicmorg/java/21"
autostart=true
autorestart=true
killasgroup=true
stopasgroup=true
#stdout_logfile=/var/log/supervisor/teamcity-agent.out.log
#stderr_logfile=/var/log/supervisor/teamcity-agent.err.log
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
user=root
priority=10
startretries=5
