# linux/ecosystem/github/runner/latest/etc/supervisor/conf.d/docker-service.conf
[program:dockerd]
command=/usr/bin/dockerd --host unix:///var/run/docker.sock --experimental --iptables=false --bridge=none --log-level=warn
autostart=true
autorestart=true
killasgroup=true
stopasgroup=true
stdout_logfile=/var/log/supervisor/dockerd.out.log
stdout_logfile_maxbytes=128MB
stderr_logfile=/var/log/supervisor/dockerd.err.log
stderr_logfile_maxbytes=128MB
user=root
priority=1
startretries=5
