# linux/ecosystem/gitlab/runner/minimal/etc/supervisor/conf.d/gitlab-runner.conf
[program:gitlab-runner]
command=/usr/local/share/epicmorg/gitlab/runner/bin/gitlab-runner run --working-directory /usr/local/share/epicmorg/gitlab/runner/worker --config /usr/local/share/epicmorg/gitlab/runner/etc/gitlab-runner/config.toml --service gitlab-runner
autostart=true
autorestart=true
killasgroup=true
stopasgroup=true
#stdout_logfile=/var/log/supervisor/gitlab-runner.out.log
#stderr_logfile=/var/log/supervisor/gitlab-runner.err.log
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
user=root
priority=10
startretries=5
